# ansible/ansible-speedtest/files/speed-test.service
[Unit]
Description=Run Speedtest CLI and log the output

[Service]
Type=oneshot
ExecStart=/usr/local/bin/speed-test.sh
